# daemon/simple_daemon.service
[Unit]
Description=Super simple daemon

[Service]
Type=simple
ExecStart=/opt/simple_daemon/simple_daemon

[Install]
WantedBy=multi-user.target
